# pilotpanel.service
[Unit]
Description=PilotPanel LXC Control Panel
After=network.target

[Service]
User=root
WorkingDirectory=/var/www/lxc
ExecStart=/var/www/lxc/venv/bin/python app.py
Restart=always
RestartSec=5
Environment=PATH=/var/www/lxc/venv/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
