# bot/mintyhost-bot.service
[Unit]
Description=MintyHost VPS Management Discord Bot
After=network.target mintyhost.service

[Service]
User=root
WorkingDirectory=/root/lxc/bot
ExecStart=/root/lxc/bot/venv/bin/python bot.py
Restart=always
RestartSec=5
Environment=PATH=/root/lxc/bot/venv/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
