# bot/pilotpanel-bot.service
[Unit]
Description=PilotPanel VPS Management Discord Bot
After=network.target pilotpanel.service

[Service]
User=root
WorkingDirectory=/var/www/lxc/bot
ExecStart=/var/www/lxc/bot/venv/bin/python bot.py
Restart=always
RestartSec=5
Environment=PATH=/var/www/lxc/bot/venv/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
